# .deploy/fedidevs-worker.service
[Unit]
Description=Fedidevs Worker
After=fedidevs.service

[Service]
Type=simple
Restart=always
RestartSec=1
WorkingDirectory=/var/apps/fedidevs
ExecStart=/var/apps/fedidevs/.venv/bin/newrelic-admin run-program /var/apps/fedidevs/.venv/bin/celery -A fedidevs worker --loglevel=info --concurrency=2
Environment="NEW_RELIC_CONFIG_FILE=/var/apps/fedidevs/newrelic.ini"
User=anze

[Install]
WantedBy=fedidevs.service
